# scripts/supervisord/supervisord-laravel.conf
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work --stop-when-empty --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/app/storage/logs/laravel-worker.log
exitcodes=0
startretries=1000
startsecs=5
stopwaitsecs=3600




[program:octane]
command=php /app/artisan octane:frankenphp --workers=300 --max-requests=1000 --host=0.0.0.0 --port=8030
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/app/storage/logs/octane.log
